# relay/com.herdr-remote.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.herdr-remote.relay</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/apple/Developer/dcolinmorgan/herdi/relay/.venv/bin/python3</string>
        <string>/Users/apple/Developer/dcolinmorgan/herdi/relay/herdr_relay.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/herdr-remote-relay.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/herdr-remote-relay.err</string>
</dict>
</plist>
